# misc/systemd/daemon.socket.in
[Unit]
Description=Lix Daemon Socket (@desc@)
Before=multi-user.target
Conflicts=nix-daemon.service
RequiresMountsFor=@storedir@
ConditionPathIsReadWrite=@localstatedir@/nix/daemon-socket

[Socket]
ListenStream=@localstatedir@/nix/daemon-socket/@path@
FileDescriptorName=@protocol@
Accept=yes

[Install]
WantedBy=sockets.target
