# misc/systemd/nix-daemon@.service.in
[Unit]
Description=Nix Daemon instance
Documentation=man:nix-daemon https://docs.lix.systems/manual/lix/stable
CollectMode=inactive-or-failed
RequiresMountsFor=@storedir@
RequiresMountsFor=@localstatedir@
RequiresMountsFor=@localstatedir@/nix/db

[Service]
ExecStart=@@bindir@/nix-daemon nix-daemon --for-socket-activation
CacheDirectory=nix
LimitNOFILE=1048576
TasksMax=1048576
Delegate=yes
DelegateSubgroup=supervisor
